# migration/.docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisor/supervisord.pid
logfile = /var/log/supervisord.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log

[program:redis]
command=redis-server /etc/redis/redis.conf
pidfile=/var/run/redis/redis.pid
logfile=/var/log/redis/redis.log
directory=/usr/bin
autostart=true
autorestart=true
redirect_stderr=true